# aws_backup.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
PYENV_VERSION=3.4.5
15 3 1 * * root /root/aws-backup/aws_backup --volume-ids vol-XXXXXXXX --mounts / --tags frequency=monthly --keep 1 --mysql --mysql-credentials ~/.aws_mysql.json
30 3 * * 0 root /root/aws-backup/aws_backup --volume-ids vol-XXXXXXXX --mounts / --tags frequency=weekly --keep 1 --mysql --mysql-credentials ~/.aws_mysql.json
45 3 * * * root /root/aws-backup/aws_backup --volume-ids vol-XXXXXXXX --mounts / --tags frequency=daily --keep 6 --mysql --mysql-credentials ~/.aws_mysql.json
